# fixtures/attach-processes.service
[Unit]
Description=attach processes test

[Service]
Delegate=yes
ExecStart=/bin/sleep 400
