# fixtures/cancelme.service
[Unit]
Description=cancel unit test

[Service]
ExecStartPre=/bin/sleep 400
ExecStart=/bin/sleep 400
