# fixtures/get-unit-pid.service
[Unit]
Description=get unit pid

[Service]
ExecStart=get-unit-pid.sh
