# fixtures/oneshot.service
[Unit]
Description=list jobs test

[Service]
Type=oneshot
ExecStart=/bin/sleep 30
